# supervisor/conf.d/crond.conf
[program:crond]
command=crond -f -L /var/log/crond.log
