# supervisor/conf.d/dhcpd.conf
[program:dhcpd]
command=/usr/share/dhcpd/dhcpd.sh -f -cf /etc/dhcpd.conf -lf /data/dhcpd.leases --no-pid
