# supervisor/conf.d/rsyslogd.conf
[program:rsyslog]
command=rsyslogd -n -f /etc/rsyslog.conf
